# scripts/watchdog/systemd/nuc5-g4-monitor.timer
[Unit]
Description=Monitor G4 every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Unit=g4-monitor.service

[Install]
WantedBy=timers.target
