# scripts/watchdog/systemd/nuc5-g4-monitor.service
[Unit]
Description=G4 compute node monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/g4-monitor.sh
